# launchd/com.syntax922.ig88.param-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.syntax922.ig88.param-proxy</string>
    <key>UserName</key>
    <string>copilot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>-u</string>
        <string>/Users/copilot/ig88-host/scripts/lmstudio-param-proxy.py</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/copilot/ig88-host/logs/param-proxy.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/copilot/ig88-host/logs/param-proxy.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/copilot/ig88-host</string>
</dict>
</plist>
