# launchd/com.syntax922.ig88.lmstudio.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.syntax922.ig88.lmstudio</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/copilot/ig88-host/scripts/lms-wrapper.sh</string>
  </array>
  <key>UserName</key>
  <string>copilot</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>WorkingDirectory</key>
  <string>/Users/copilot/ig88-host</string>
  <key>StandardOutPath</key>
  <string>/Users/copilot/ig88-host/logs/lmstudio.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/copilot/ig88-host/logs/lmstudio.log</string>
</dict>
</plist>
